# puffer-update.service
[Unit]
Description=Pull latest lifesupport code and restart display if changed

[Service]
Type=oneshot
WorkingDirectory=/home/pi/aquarium-life-support
ExecStart=/bin/bash -c '\
  before=$(git -c safe.directory=/home/pi/aquarium-life-support rev-parse HEAD) && \
  git -c safe.directory=/home/pi/aquarium-life-support pull --ff-only && \
  after=$(git -c safe.directory=/home/pi/aquarium-life-support rev-parse HEAD) && \
  if [ "$before" != "$after" ]; then \
    echo "New commits detected, restarting puffer-display..." && \
    systemctl restart puffer-display.service; \
  else \
    echo "Already up to date."; \
  fi'
